# medical-dictation-backend/infra/hf-space/supervisord.conf
; Sprint-07 HF Space supervisord.
;
; Auth model: real Keycloak (sprint-02 realm), seeded test clinician.
; HF API key gates access to the Space itself.

[supervisord]
nodaemon=true
user=root
logfile=/var/log/medical-dictation/supervisord.log
loglevel=info

[program:postgres]
command=/usr/lib/postgresql/16/bin/postgres -D /run/pgdata -c listen_addresses=127.0.0.1
user=postgres
autostart=true
autorestart=false
stdout_logfile=/var/log/medical-dictation/postgres.log
stderr_logfile=/var/log/medical-dictation/postgres.err
priority=10

[program:redis]
command=/usr/bin/redis-server --bind 127.0.0.1 --save '' --appendonly no
user=redis
autostart=true
autorestart=true
stdout_logfile=/var/log/medical-dictation/redis.log
stderr_logfile=/var/log/medical-dictation/redis.err
priority=10

[program:keycloak]
; start-dev because the import-realm flow + tmpfs DB makes production mode
; awkward in a single-container demo. ADR-0017 documents.
command=/opt/keycloak/bin/kc.sh start-dev --import-realm --http-port=8088
environment=
    KC_DB=postgres,
    KC_DB_URL=jdbc:postgresql://127.0.0.1:5432/keycloak,
    KC_DB_USERNAME=keycloak,
    KC_DB_PASSWORD=keycloak,
    KC_HTTP_RELATIVE_PATH=/auth,
    KC_HOSTNAME_STRICT=false,
    KC_HOSTNAME_STRICT_HTTPS=false,
    KEYCLOAK_ADMIN=admin,
    KEYCLOAK_ADMIN_PASSWORD=admin
autostart=true
autorestart=true
priority=15
stdout_logfile=/var/log/medical-dictation/keycloak.log
stderr_logfile=/var/log/medical-dictation/keycloak.err

[program:auth-service]
command=/usr/bin/python3.12 -m uvicorn auth_service.main:app --host 127.0.0.1 --port 8000
environment=
    ENVIRONMENT=demo,
    AUTH_ISSUER=http://127.0.0.1:8088/auth/realms/medical-dictation,
    AUTH_JWKS_URL=http://127.0.0.1:8088/auth/realms/medical-dictation/protocol/openid-connect/certs,
    KEYCLOAK_BASE_URL=http://127.0.0.1:8088/auth,
    DB_APP_ROLE_DSN=postgresql://app_role:app_role@127.0.0.1:5432/medical_dictation,
    DB_AUDIT_WRITER_DSN=postgresql://audit_writer:audit_writer@127.0.0.1:5432/medical_dictation,
    DB_TENANT_WRITER_DSN=postgresql://tenant_writer:tenant_writer@127.0.0.1:5432/medical_dictation,
    REDIS_URL=redis://127.0.0.1:6379/0
autostart=true
autorestart=true
priority=20
stdout_logfile=/var/log/medical-dictation/auth-service.log
stderr_logfile=/var/log/medical-dictation/auth-service.err

[program:report-service]
command=/usr/bin/python3.12 -m uvicorn report_service.main:app --host 127.0.0.1 --port 8006
environment=
    ENVIRONMENT=demo,
    AUTH_ISSUER=http://127.0.0.1:8088/auth/realms/medical-dictation,
    AUTH_JWKS_URL=http://127.0.0.1:8088/auth/realms/medical-dictation/protocol/openid-connect/certs,
    DB_APP_ROLE_DSN=postgresql://app_role:app_role@127.0.0.1:5432/medical_dictation,
    DB_AUDIT_WRITER_DSN=postgresql://audit_writer:audit_writer@127.0.0.1:5432/medical_dictation
autostart=true
autorestart=true
priority=21
stdout_logfile=/var/log/medical-dictation/report-service.log
stderr_logfile=/var/log/medical-dictation/report-service.err

[program:nlp-service]
command=/usr/bin/python3.12 -m uvicorn nlp_service.main:app --host 127.0.0.1 --port 8005
environment=
    ENVIRONMENT=demo,
    AUTH_ISSUER=http://127.0.0.1:8088/auth/realms/medical-dictation,
    AUTH_JWKS_URL=http://127.0.0.1:8088/auth/realms/medical-dictation/protocol/openid-connect/certs,
    DB_APP_ROLE_DSN=postgresql://app_role:app_role@127.0.0.1:5432/medical_dictation,
    DB_AUDIT_WRITER_DSN=postgresql://audit_writer:audit_writer@127.0.0.1:5432/medical_dictation,
    REDIS_URL=redis://127.0.0.1:6379/0,
    MDX_NLP_PUNCTUATION_DISABLED=true
autostart=true
autorestart=true
priority=22
stdout_logfile=/var/log/medical-dictation/nlp-service.log
stderr_logfile=/var/log/medical-dictation/nlp-service.err

[program:asr-service]
command=/usr/bin/python3.12 -m uvicorn asr_service.main:app --host 127.0.0.1 --port 8001
environment=
    ENVIRONMENT=demo,
    MD_OBJECT_STORE_DISABLED=true,
    AUTH_ISSUER=http://127.0.0.1:8088/auth/realms/medical-dictation,
    AUTH_JWKS_URL=http://127.0.0.1:8088/auth/realms/medical-dictation/protocol/openid-connect/certs,
    DB_APP_ROLE_DSN=postgresql://app_role:app_role@127.0.0.1:5432/medical_dictation,
    REDIS_URL=redis://127.0.0.1:6379/0
autostart=true
autorestart=true
priority=23
stdout_logfile=/var/log/medical-dictation/asr-service.log
stderr_logfile=/var/log/medical-dictation/asr-service.err

[program:asr-worker]
command=/usr/bin/python3.12 -m asr_worker.main
environment=
    ENVIRONMENT=demo,
    MD_OBJECT_STORE_DISABLED=true,
    MD_ASR_DEVICE=cuda,
    MD_ASR_MODEL=large-v3,
    DB_APP_ROLE_DSN=postgresql://app_role:app_role@127.0.0.1:5432/medical_dictation,
    REDIS_URL=redis://127.0.0.1:6379/0
autostart=true
autorestart=true
priority=24
stdout_logfile=/var/log/medical-dictation/asr-worker.log
stderr_logfile=/var/log/medical-dictation/asr-worker.err

[program:dictation-service]
command=/usr/bin/python3.12 -m uvicorn dictation_service.main:app --host 127.0.0.1 --port 8002
environment=
    ENVIRONMENT=demo,
    MD_OBJECT_STORE_DISABLED=true,
    DEMO_AUDIO_PURGE_ON_FINALIZE=true,
    MDX_TMPFS_ROOT=/run/dictation,
    AUTH_ISSUER=http://127.0.0.1:8088/auth/realms/medical-dictation,
    AUTH_JWKS_URL=http://127.0.0.1:8088/auth/realms/medical-dictation/protocol/openid-connect/certs,
    DB_APP_ROLE_DSN=postgresql://app_role:app_role@127.0.0.1:5432/medical_dictation,
    REDIS_URL=redis://127.0.0.1:6379/0
autostart=true
autorestart=true
priority=25
stdout_logfile=/var/log/medical-dictation/dictation-service.log
stderr_logfile=/var/log/medical-dictation/dictation-service.err

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
autostart=true
autorestart=true
priority=50
stdout_logfile=/var/log/medical-dictation/nginx.log
stderr_logfile=/var/log/medical-dictation/nginx.err

[program:daily-restart-cron]
; Force a container restart at 06:00 UTC daily (privacy hygiene).
command=/usr/bin/python3.12 -c "import time, os, signal; t=time.gmtime(); secs=((6-t.tm_hour) % 24)*3600 - t.tm_min*60 - t.tm_sec; secs=secs+86400 if secs<=0 else secs; time.sleep(secs); os.kill(1, signal.SIGTERM)"
autostart=true
autorestart=true
priority=99
stdout_logfile=/var/log/medical-dictation/daily-restart.log
stderr_logfile=/var/log/medical-dictation/daily-restart.err
